# fedzk/examples/service-scripts/systemd/fedzk-coordinator.service
[Unit]
Description=FedZK Coordinator Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=fedzk
Group=fedzk
WorkingDirectory=/opt/fedzk
ExecStart=/opt/fedzk/examples/service-scripts/run-coordinator.sh --port 8001 --host 0.0.0.0
Restart=on-failure
RestartSec=5

# Security hardening
ProtectSystem=strict
ReadWritePaths=/opt/fedzk
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target 