# fedzk/examples/service-scripts/systemd/fedzk-mpc.service
[Unit]
Description=FedZK MPC Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=fedzk
Group=fedzk
WorkingDirectory=/opt/fedzk
ExecStart=/opt/fedzk/examples/service-scripts/run-mpc-server.sh --port 8000 --host 0.0.0.0
Restart=on-failure
RestartSec=5
Environment="MPC_API_KEYS=%API_KEYS%"

# Security hardening
ProtectSystem=strict
ReadWritePaths=/opt/fedzk
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target 